# lib/dpkg/info/alsa-utils.conffiles
/etc/init.d/alsa-utils
/etc/init/alsa-utils.conf
