# lib/dpkg/info/avahi-daemon.conffiles
/etc/avahi/avahi-daemon.conf
/etc/avahi/hosts
/etc/dbus-1/system.d/avahi-dbus.conf
/etc/default/avahi-daemon
/etc/init.d/avahi-daemon
/etc/init/avahi-cups-reload.conf
/etc/init/avahi-daemon.conf
/etc/network/if-up.d/avahi-daemon
/etc/resolvconf/update-libc.d/avahi-daemon
