# lib/dpkg/info/binfmt-support.conffiles
/etc/init.d/binfmt-support
/etc/init/binfmt-support.conf
