# lib/dpkg/info/binfmt-support.list
/.
/etc
/etc/init.d
/etc/init.d/binfmt-support
/etc/init
/etc/init/binfmt-support.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/binfmt-support.service
/var
/var/lib
/var/lib/binfmts
/usr
/usr/sbin
/usr/sbin/update-binfmts
/usr/lib
/usr/lib/binfmt-support
/usr/lib/binfmt-support/run-detectors
/usr/share
/usr/share/doc
/usr/share/doc/binfmt-support
/usr/share/doc/binfmt-support/copyright
/usr/share/doc/binfmt-support/changelog.Debian.gz
/usr/share/doc/binfmt-support/README.Debian
/usr/share/doc/binfmt-support/detectors
/usr/share/doc/binfmt-support/TODO.Debian
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/update-binfmts.8.gz
/usr/share/binfmts
