# lib/dpkg/info/bluez.list
/.
/bin
/bin/hciconfig
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/bluetooth.service
/lib/udev
/lib/udev/hid2hci
/lib/udev/rules.d
/lib/udev/rules.d/50-bluetooth-hci-auto-poweron.rules
/lib/udev/rules.d/97-hid2hci.rules
/usr
/usr/bin
/usr/bin/bccmd
/usr/bin/hciattach
/usr/bin/bluetoothctl
/usr/bin/gatttool
/usr/bin/bluemoon
/usr/bin/rctest
/usr/bin/hex2hcd
/usr/bin/obexctl
/usr/bin/btmgmt
/usr/bin/sdptool
/usr/bin/l2test
/usr/bin/ciptool
/usr/bin/rfcomm
/usr/bin/hcitool
/usr/bin/l2ping
/usr/bin/btmon
/usr/bin/btattach
/usr/sbin
/usr/lib
/usr/lib/bluetooth
/usr/lib/bluetooth/bluetoothd
/usr/lib/bluetooth/mpris-proxy
/usr/share
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.bluez.service
/usr/share/doc
/usr/share/doc/bluez
/usr/share/doc/bluez/NEWS.Debian.gz
/usr/share/doc/bluez/copyright
/usr/share/doc/bluez/README.Debian.gz
/usr/share/doc/bluez/changelog.Debian.gz
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/sdptool.1.gz
/usr/share/man/man1/hcitool.1.gz
/usr/share/man/man1/hciattach.1.gz
/usr/share/man/man1/btmon.1.gz
/usr/share/man/man1/bluetoothctl.1.gz
/usr/share/man/man1/l2ping.1.gz
/usr/share/man/man1/hid2hci.1.gz
/usr/share/man/man1/btattach.1.gz
/usr/share/man/man1/hciconfig.1.gz
/usr/share/man/man1/mpris-proxy.1.gz
/usr/share/man/man1/bccmd.1.gz
/usr/share/man/man1/bluemoon.1.gz
/usr/share/man/man1/hex2hcd.1.gz
/usr/share/man/man1/rfcomm.1.gz
/usr/share/man/man1/gatttool.1.gz
/usr/share/man/man1/ciptool.1.gz
/usr/share/man/man1/l2test.1.gz
/usr/share/man/man1/rctest.1.gz
/usr/share/man/man8
/usr/share/man/man8/bluetoothd.8.gz
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_bluez.py
/etc
/etc/init
/etc/init/bluetooth.conf
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/bluetooth.conf
/etc/bluetooth
/etc/bluetooth/main.conf
/etc/bluetooth/proximity.conf
/etc/bluetooth/input.conf
/etc/bluetooth/network.conf
/etc/init.d
/etc/init.d/bluetooth
/var
/var/lib
/var/lib/bluetooth
/usr/sbin/bluetoothd
