# lib/dpkg/info/dbus.conffiles
/etc/default/dbus
/etc/init.d/dbus
/etc/init/dbus.conf
