# lib/dpkg/info/gdm3.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/gdm.service
/usr
/usr/sbin
/usr/sbin/gdm3
/usr/lib
/usr/lib/gdm3
/usr/lib/gdm3/gdm-session-worker
/usr/lib/gdm3/gdm-wayland-session
/usr/lib/gdm3/gdm-host-chooser
/usr/lib/gdm3/gdm-x-session
/usr/lib/gdm3/gdm-simple-chooser
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/gdm3
/usr/share/help
/usr/share/help/C
/usr/share/help/C/gdm
/usr/share/help/C/gdm/legal.xml
/usr/share/help/C/gdm/index.docbook
/usr/share/help/id
/usr/share/help/id/gdm
/usr/share/help/id/gdm/legal.xml
/usr/share/help/id/gdm/index.docbook
/usr/share/help/ca
/usr/share/help/ca/gdm
/usr/share/help/ca/gdm/legal.xml
/usr/share/help/ca/gdm/index.docbook
/usr/share/help/sl
/usr/share/help/sl/gdm
/usr/share/help/sl/gdm/legal.xml
/usr/share/help/sl/gdm/index.docbook
/usr/share/help/uk
/usr/share/help/uk/gdm
/usr/share/help/uk/gdm/legal.xml
/usr/share/help/uk/gdm/index.docbook
/usr/share/help/en_GB
/usr/share/help/en_GB/gdm
/usr/share/help/en_GB/gdm/legal.xml
/usr/share/help/en_GB/gdm/index.docbook
/usr/share/help/es
/usr/share/help/es/gdm
/usr/share/help/es/gdm/legal.xml
/usr/share/help/es/gdm/index.docbook
/usr/share/help/ro
/usr/share/help/ro/gdm
/usr/share/help/ro/gdm/legal.xml
/usr/share/help/ro/gdm/index.docbook
/usr/share/help/ru
/usr/share/help/ru/gdm
/usr/share/help/ru/gdm/legal.xml
/usr/share/help/ru/gdm/index.docbook
/usr/share/help/zh_CN
/usr/share/help/zh_CN/gdm
/usr/share/help/zh_CN/gdm/legal.xml
/usr/share/help/zh_CN/gdm/index.docbook
/usr/share/help/pt_BR
/usr/share/help/pt_BR/gdm
/usr/share/help/pt_BR/gdm/legal.xml
/usr/share/help/pt_BR/gdm/index.docbook
/usr/share/help/gl
/usr/share/help/gl/gdm
/usr/share/help/gl/gdm/legal.xml
/usr/share/help/gl/gdm/index.docbook
/usr/share/help/fr
/usr/share/help/fr/gdm
/usr/share/help/fr/gdm/legal.xml
/usr/share/help/fr/gdm/index.docbook
/usr/share/help/sv
/usr/share/help/sv/gdm
/usr/share/help/sv/gdm/legal.xml
/usr/share/help/sv/gdm/index.docbook
/usr/share/help/oc
/usr/share/help/oc/gdm
/usr/share/help/oc/gdm/legal.xml
/usr/share/help/oc/gdm/index.docbook
/usr/share/help/hu
/usr/share/help/hu/gdm
/usr/share/help/hu/gdm/legal.xml
/usr/share/help/hu/gdm/index.docbook
/usr/share/help/it
/usr/share/help/it/gdm
/usr/share/help/it/gdm/legal.xml
/usr/share/help/it/gdm/index.docbook
/usr/share/help/el
/usr/share/help/el/gdm
/usr/share/help/el/gdm/legal.xml
/usr/share/help/el/gdm/index.docbook
/usr/share/help/de
/usr/share/help/de/gdm
/usr/share/help/de/gdm/legal.xml
/usr/share/help/de/gdm/index.docbook
/usr/share/help/te
/usr/share/help/te/gdm
/usr/share/help/te/gdm/legal.xml
/usr/share/help/te/gdm/index.docbook
/usr/share/help/cs
/usr/share/help/cs/gdm
/usr/share/help/cs/gdm/legal.xml
/usr/share/help/cs/gdm/index.docbook
/usr/share/help/ko
/usr/share/help/ko/gdm
/usr/share/help/ko/gdm/legal.xml
/usr/share/help/ko/gdm/index.docbook
/usr/share/pixmaps
/usr/share/pixmaps/nohost.png
/usr/share/pixmaps/nobody.png
/usr/share/pixmaps/gdm.png
/usr/share/pixmaps/gdm-setup.png
/usr/share/pixmaps/gdm-xnest.png
/usr/share/pixmaps/gdm-foot-logo.png
/usr/share/dconf
/usr/share/dconf/profile
/usr/share/dconf/profile/gdm
/usr/share/locale
/usr/share/locale/rw
/usr/share/locale/rw/LC_MESSAGES
/usr/share/locale/rw/LC_MESSAGES/gdm.mo
/usr/share/locale/ta
/usr/share/locale/ta/LC_MESSAGES
/usr/share/locale/ta/LC_MESSAGES/gdm.mo
/usr/share/locale/mk
/usr/share/locale/mk/LC_MESSAGES
/usr/share/locale/mk/LC_MESSAGES/gdm.mo
/usr/share/locale/sq
/usr/share/locale/sq/LC_MESSAGES
/usr/share/locale/sq/LC_MESSAGES/gdm.mo
/usr/share/locale/id
/usr/share/locale/id/LC_MESSAGES
/usr/share/locale/id/LC_MESSAGES/gdm.mo
/usr/share/locale/pa
/usr/share/locale/pa/LC_MESSAGES
/usr/share/locale/pa/LC_MESSAGES/gdm.mo
/usr/share/locale/ca
/usr/share/locale/ca/LC_MESSAGES
/usr/share/locale/ca/LC_MESSAGES/gdm.mo
/usr/share/locale/km
/usr/share/locale/km/LC_MESSAGES
/usr/share/locale/km/LC_MESSAGES/gdm.mo
/usr/share/locale/sr
/usr/share/locale/sr/LC_MESSAGES
/usr/share/locale/sr/LC_MESSAGES/gdm.mo
/usr/share/locale/crh
/usr/share/locale/crh/LC_MESSAGES
/usr/share/locale/crh/LC_MESSAGES/gdm.mo
/usr/share/locale/bg
/usr/share/locale/bg/LC_MESSAGES
/usr/share/locale/bg/LC_MESSAGES/gdm.mo
/usr/share/locale/ja
/usr/share/locale/ja/LC_MESSAGES
/usr/share/locale/ja/LC_MESSAGES/gdm.mo
/usr/share/locale/nso
/usr/share/locale/nso/LC_MESSAGES
/usr/share/locale/nso/LC_MESSAGES/gdm.mo
/usr/share/locale/fur
/usr/share/locale/fur/LC_MESSAGES
/usr/share/locale/fur/LC_MESSAGES/gdm.mo
/usr/share/locale/sr@latin
/usr/share/locale/sr@latin/LC_MESSAGES
/usr/share/locale/sr@latin/LC_MESSAGES/gdm.mo
/usr/share/locale/da
/usr/share/locale/da/LC_MESSAGES
/usr/share/locale/da/LC_MESSAGES/gdm.mo
/usr/share/locale/fa
/usr/share/locale/fa/LC_MESSAGES
/usr/share/locale/fa/LC_MESSAGES/gdm.mo
/usr/share/locale/wa
/usr/share/locale/wa/LC_MESSAGES
/usr/share/locale/wa/LC_MESSAGES/gdm.mo
/usr/share/locale/am
/usr/share/locale/am/LC_MESSAGES
/usr/share/locale/am/LC_MESSAGES/gdm.mo
/usr/share/locale/an
/usr/share/locale/an/LC_MESSAGES
/usr/share/locale/an/LC_MESSAGES/gdm.mo
/usr/share/locale/as
/usr/share/locale/as/LC_MESSAGES
/usr/share/locale/as/LC_MESSAGES/gdm.mo
/usr/share/locale/fi
/usr/share/locale/fi/LC_MESSAGES
/usr/share/locale/fi/LC_MESSAGES/gdm.mo
/usr/share/locale/ka
/usr/share/locale/ka/LC_MESSAGES
/usr/share/locale/ka/LC_MESSAGES/gdm.mo
/usr/share/locale/zu
/usr/share/locale/zu/LC_MESSAGES
/usr/share/locale/zu/LC_MESSAGES/gdm.mo
/usr/share/locale/sl
/usr/share/locale/sl/LC_MESSAGES
/usr/share/locale/sl/LC_MESSAGES/gdm.mo
/usr/share/locale/ar
/usr/share/locale/ar/LC_MESSAGES
/usr/share/locale/ar/LC_MESSAGES/gdm.mo
/usr/share/locale/uk
/usr/share/locale/uk/LC_MESSAGES
/usr/share/locale/uk/LC_MESSAGES/gdm.mo
/usr/share/locale/bs
/usr/share/locale/bs/LC_MESSAGES
/usr/share/locale/bs/LC_MESSAGES/gdm.mo
/usr/share/locale/bn
/usr/share/locale/bn/LC_MESSAGES
/usr/share/locale/bn/LC_MESSAGES/gdm.mo
/usr/share/locale/en_GB
/usr/share/locale/en_GB/LC_MESSAGES
/usr/share/locale/en_GB/LC_MESSAGES/gdm.mo
/usr/share/locale/he
/usr/share/locale/he/LC_MESSAGES
/usr/share/locale/he/LC_MESSAGES/gdm.mo
/usr/share/locale/ga
/usr/share/locale/ga/LC_MESSAGES
/usr/share/locale/ga/LC_MESSAGES/gdm.mo
/usr/share/locale/es
/usr/share/locale/es/LC_MESSAGES
/usr/share/locale/es/LC_MESSAGES/gdm.mo
/usr/share/locale/ro
/usr/share/locale/ro/LC_MESSAGES
/usr/share/locale/ro/LC_MESSAGES/gdm.mo
/usr/share/locale/bn_IN
/usr/share/locale/bn_IN/LC_MESSAGES
/usr/share/locale/bn_IN/LC_MESSAGES/gdm.mo
/usr/share/locale/ru
/usr/share/locale/ru/LC_MESSAGES
/usr/share/locale/ru/LC_MESSAGES/gdm.mo
/usr/share/locale/th
/usr/share/locale/th/LC_MESSAGES
/usr/share/locale/th/LC_MESSAGES/gdm.mo
/usr/share/locale/zh_CN
/usr/share/locale/zh_CN/LC_MESSAGES
/usr/share/locale/zh_CN/LC_MESSAGES/gdm.mo
/usr/share/locale/tr
/usr/share/locale/tr/LC_MESSAGES
/usr/share/locale/tr/LC_MESSAGES/gdm.mo
/usr/share/locale/ne
/usr/share/locale/ne/LC_MESSAGES
/usr/share/locale/ne/LC_MESSAGES/gdm.mo
/usr/share/locale/ca@valencia
/usr/share/locale/ca@valencia/LC_MESSAGES
/usr/share/locale/ca@valencia/LC_MESSAGES/gdm.mo
/usr/share/locale/pt_BR
/usr/share/locale/pt_BR/LC_MESSAGES
/usr/share/locale/pt_BR/LC_MESSAGES/gdm.mo
/usr/share/locale/ml
/usr/share/locale/ml/LC_MESSAGES
/usr/share/locale/ml/LC_MESSAGES/gdm.mo
/usr/share/locale/gl
/usr/share/locale/gl/LC_MESSAGES
/usr/share/locale/gl/LC_MESSAGES/gdm.mo
/usr/share/locale/fr
/usr/share/locale/fr/LC_MESSAGES
/usr/share/locale/fr/LC_MESSAGES/gdm.mo
/usr/share/locale/br
/usr/share/locale/br/LC_MESSAGES
/usr/share/locale/br/LC_MESSAGES/gdm.mo
/usr/share/locale/ms
/usr/share/locale/ms/LC_MESSAGES
/usr/share/locale/ms/LC_MESSAGES/gdm.mo
/usr/share/locale/af
/usr/share/locale/af/LC_MESSAGES
/usr/share/locale/af/LC_MESSAGES/gdm.mo
/usr/share/locale/lt
/usr/share/locale/lt/LC_MESSAGES
/usr/share/locale/lt/LC_MESSAGES/gdm.mo
/usr/share/locale/nb
/usr/share/locale/nb/LC_MESSAGES
/usr/share/locale/nb/LC_MESSAGES/gdm.mo
/usr/share/locale/kab
/usr/share/locale/kab/LC_MESSAGES
/usr/share/locale/kab/LC_MESSAGES/gdm.mo
/usr/share/locale/ku
/usr/share/locale/ku/LC_MESSAGES
/usr/share/locale/ku/LC_MESSAGES/gdm.mo
/usr/share/locale/kn
/usr/share/locale/kn/LC_MESSAGES
/usr/share/locale/kn/LC_MESSAGES/gdm.mo
/usr/share/locale/gu
/usr/share/locale/gu/LC_MESSAGES
/usr/share/locale/gu/LC_MESSAGES/gdm.mo
/usr/share/locale/fy
/usr/share/locale/fy/LC_MESSAGES
/usr/share/locale/fy/LC_MESSAGES/gdm.mo
/usr/share/locale/be
/usr/share/locale/be/LC_MESSAGES
/usr/share/locale/be/LC_MESSAGES/gdm.mo
/usr/share/locale/ast
/usr/share/locale/ast/LC_MESSAGES
/usr/share/locale/ast/LC_MESSAGES/gdm.mo
/usr/share/locale/ky
/usr/share/locale/ky/LC_MESSAGES
/usr/share/locale/ky/LC_MESSAGES/gdm.mo
/usr/share/locale/mr
/usr/share/locale/mr/LC_MESSAGES
/usr/share/locale/mr/LC_MESSAGES/gdm.mo
/usr/share/locale/mg
/usr/share/locale/mg/LC_MESSAGES
/usr/share/locale/mg/LC_MESSAGES/gdm.mo
/usr/share/locale/nds
/usr/share/locale/nds/LC_MESSAGES
/usr/share/locale/nds/LC_MESSAGES/gdm.mo
/usr/share/locale/dz
/usr/share/locale/dz/LC_MESSAGES
/usr/share/locale/dz/LC_MESSAGES/gdm.mo
/usr/share/locale/eu
/usr/share/locale/eu/LC_MESSAGES
/usr/share/locale/eu/LC_MESSAGES/gdm.mo
/usr/share/locale/nl
/usr/share/locale/nl/LC_MESSAGES
/usr/share/locale/nl/LC_MESSAGES/gdm.mo
/usr/share/locale/sv
/usr/share/locale/sv/LC_MESSAGES
/usr/share/locale/sv/LC_MESSAGES/gdm.mo
/usr/share/locale/oc
/usr/share/locale/oc/LC_MESSAGES
/usr/share/locale/oc/LC_MESSAGES/gdm.mo
/usr/share/locale/hu
/usr/share/locale/hu/LC_MESSAGES
/usr/share/locale/hu/LC_MESSAGES/gdm.mo
/usr/share/locale/az
/usr/share/locale/az/LC_MESSAGES
/usr/share/locale/az/LC_MESSAGES/gdm.mo
/usr/share/locale/nn
/usr/share/locale/nn/LC_MESSAGES
/usr/share/locale/nn/LC_MESSAGES/gdm.mo
/usr/share/locale/zh_HK
/usr/share/locale/zh_HK/LC_MESSAGES
/usr/share/locale/zh_HK/LC_MESSAGES/gdm.mo
/usr/share/locale/sk
/usr/share/locale/sk/LC_MESSAGES
/usr/share/locale/sk/LC_MESSAGES/gdm.mo
/usr/share/locale/xh
/usr/share/locale/xh/LC_MESSAGES
/usr/share/locale/xh/LC_MESSAGES/gdm.mo
/usr/share/locale/zh_TW
/usr/share/locale/zh_TW/LC_MESSAGES
/usr/share/locale/zh_TW/LC_MESSAGES/gdm.mo
/usr/share/locale/pl
/usr/share/locale/pl/LC_MESSAGES
/usr/share/locale/pl/LC_MESSAGES/gdm.mo
/usr/share/locale/mi
/usr/share/locale/mi/LC_MESSAGES
/usr/share/locale/mi/LC_MESSAGES/gdm.mo
/usr/share/locale/en_CA
/usr/share/locale/en_CA/LC_MESSAGES
/usr/share/locale/en_CA/LC_MESSAGES/gdm.mo
/usr/share/locale/mn
/usr/share/locale/mn/LC_MESSAGES
/usr/share/locale/mn/LC_MESSAGES/gdm.mo
/usr/share/locale/si
/usr/share/locale/si/LC_MESSAGES
/usr/share/locale/si/LC_MESSAGES/gdm.mo
/usr/share/locale/eo
/usr/share/locale/eo/LC_MESSAGES
/usr/share/locale/eo/LC_MESSAGES/gdm.mo
/usr/share/locale/en@shaw
/usr/share/locale/en@shaw/LC_MESSAGES
/usr/share/locale/en@shaw/LC_MESSAGES/gdm.mo
/usr/share/locale/gd
/usr/share/locale/gd/LC_MESSAGES
/usr/share/locale/gd/LC_MESSAGES/gdm.mo
/usr/share/locale/cy
/usr/share/locale/cy/LC_MESSAGES
/usr/share/locale/cy/LC_MESSAGES/gdm.mo
/usr/share/locale/tg
/usr/share/locale/tg/LC_MESSAGES
/usr/share/locale/tg/LC_MESSAGES/gdm.mo
/usr/share/locale/hr
/usr/share/locale/hr/LC_MESSAGES
/usr/share/locale/hr/LC_MESSAGES/gdm.mo
/usr/share/locale/csb
/usr/share/locale/csb/LC_MESSAGES
/usr/share/locale/csb/LC_MESSAGES/gdm.mo
/usr/share/locale/it
/usr/share/locale/it/LC_MESSAGES
/usr/share/locale/it/LC_MESSAGES/gdm.mo
/usr/share/locale/et
/usr/share/locale/et/LC_MESSAGES
/usr/share/locale/et/LC_MESSAGES/gdm.mo
/usr/share/locale/el
/usr/share/locale/el/LC_MESSAGES
/usr/share/locale/el/LC_MESSAGES/gdm.mo
/usr/share/locale/vi
/usr/share/locale/vi/LC_MESSAGES
/usr/share/locale/vi/LC_MESSAGES/gdm.mo
/usr/share/locale/hi
/usr/share/locale/hi/LC_MESSAGES
/usr/share/locale/hi/LC_MESSAGES/gdm.mo
/usr/share/locale/de
/usr/share/locale/de/LC_MESSAGES
/usr/share/locale/de/LC_MESSAGES/gdm.mo
/usr/share/locale/te
/usr/share/locale/te/LC_MESSAGES
/usr/share/locale/te/LC_MESSAGES/gdm.mo
/usr/share/locale/uz@cyrillic
/usr/share/locale/uz@cyrillic/LC_MESSAGES
/usr/share/locale/uz@cyrillic/LC_MESSAGES/gdm.mo
/usr/share/locale/pt
/usr/share/locale/pt/LC_MESSAGES
/usr/share/locale/pt/LC_MESSAGES/gdm.mo
/usr/share/locale/be@latin
/usr/share/locale/be@latin/LC_MESSAGES
/usr/share/locale/be@latin/LC_MESSAGES/gdm.mo
/usr/share/locale/kk
/usr/share/locale/kk/LC_MESSAGES
/usr/share/locale/kk/LC_MESSAGES/gdm.mo
/usr/share/locale/cs
/usr/share/locale/cs/LC_MESSAGES
/usr/share/locale/cs/LC_MESSAGES/gdm.mo
/usr/share/locale/lv
/usr/share/locale/lv/LC_MESSAGES
/usr/share/locale/lv/LC_MESSAGES/gdm.mo
/usr/share/locale/ps
/usr/share/locale/ps/LC_MESSAGES
/usr/share/locale/ps/LC_MESSAGES/gdm.mo
/usr/share/locale/gv
/usr/share/locale/gv/LC_MESSAGES
/usr/share/locale/gv/LC_MESSAGES/gdm.mo
/usr/share/locale/or
/usr/share/locale/or/LC_MESSAGES
/usr/share/locale/or/LC_MESSAGES/gdm.mo
/usr/share/locale/mai
/usr/share/locale/mai/LC_MESSAGES
/usr/share/locale/mai/LC_MESSAGES/gdm.mo
/usr/share/locale/ko
/usr/share/locale/ko/LC_MESSAGES
/usr/share/locale/ko/LC_MESSAGES/gdm.mo
/usr/share/locale/uz
/usr/share/locale/uz/LC_MESSAGES
/usr/share/locale/uz/LC_MESSAGES/gdm.mo
/usr/share/locale/hy
/usr/share/locale/hy/LC_MESSAGES
/usr/share/locale/hy/LC_MESSAGES/gdm.mo
/usr/share/locale/ug
/usr/share/locale/ug/LC_MESSAGES
/usr/share/locale/ug/LC_MESSAGES/gdm.mo
/usr/share/locale/is
/usr/share/locale/is/LC_MESSAGES
/usr/share/locale/is/LC_MESSAGES/gdm.mo
/usr/share/icons
/usr/share/icons/hicolor
/usr/share/icons/hicolor/16x16
/usr/share/icons/hicolor/16x16/apps
/usr/share/icons/hicolor/16x16/apps/gdm-xnest.png
/usr/share/icons/hicolor/32x32
/usr/share/icons/hicolor/32x32/apps
/usr/share/icons/hicolor/32x32/apps/gdm-setup.png
/usr/share/icons/hicolor/32x32/apps/gdm-xnest.png
/usr/share/gdm
/usr/share/gdm/BuiltInSessions
/usr/share/gdm/BuiltInSessions/default.desktop
/usr/share/gdm/greeter.dconf-defaults.md5sum
/usr/share/gdm/gdb-cmd
/usr/share/gdm/greeter.dconf-defaults
/usr/share/gdm/generate-config
/usr/share/gdm/gdm.schemas
/usr/share/gdm/locale.alias
/usr/share/gdm/dconf
/usr/share/gdm/dconf/00-upstream-settings
/usr/share/gdm/dconf/locks
/usr/share/gdm/dconf/locks/00-upstream-settings-locks
/usr/share/gdm/greeter
/usr/share/gdm/greeter/applications
/usr/share/gdm/greeter/applications/mime-dummy-handler.desktop
/usr/share/gdm/greeter/applications/mimeapps.list
/usr/share/gdm/greeter/applications/gnome-shell-wayland.desktop
/usr/share/gdm/greeter/applications/gnome-shell.desktop
/usr/share/gdm/greeter/autostart
/usr/share/gdm/greeter/autostart/orca-autostart.desktop
/usr/share/gdm/greeter/dbus-1
/usr/share/gdm/greeter/dbus-1/services
/usr/share/gdm/greeter/dbus-1/services/org.gnome.OnlineAccounts.service
/usr/share/gdm/greeter/dbus-1/services/org.freedesktop.Telepathy.AccountManager.service
/usr/share/gdm/greeter/dbus-1/services/org.freedesktop.Telepathy.MissionControl5.service
/usr/share/doc
/usr/share/doc/gdm3
/usr/share/doc/gdm3/copyright
/usr/share/doc/gdm3/README.Debian
/usr/share/doc/gdm3/TODO.Debian
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/gdm-screenshot.1.gz
/usr/share/man/man8
/usr/share/man/man8/gdm3.8.gz
/usr/bin
/usr/bin/gdm-screenshot
/var
/var/log
/var/log/gdm3
/var/lib
/var/lib/gdm3
/var/cache
/var/cache/gdm
/etc
/etc/gdm3
/etc/gdm3/PostLogin
/etc/gdm3/PostLogin/Default.sample
/etc/gdm3/Prime
/etc/gdm3/Prime/Default
/etc/gdm3/PostSession
/etc/gdm3/PostSession/Default
/etc/gdm3/Xsession
/etc/gdm3/Init
/etc/gdm3/Init/Default
/etc/gdm3/PrimeOff
/etc/gdm3/PrimeOff/Default
/etc/gdm3/PreSession
/etc/gdm3/PreSession/Default
/etc/gdm3/config-error-dialog.sh
/etc/gdm3/custom.conf
/etc/init.d
/etc/init.d/gdm3
/etc/pam.d
/etc/pam.d/gdm-autologin
/etc/pam.d/gdm-launch-environment
/etc/pam.d/gdm-password
/etc/init
/etc/init/gdm3.conf
/etc/insserv.conf.d
/etc/insserv.conf.d/gdm3
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/gdm.conf
/lib/systemd/system/gdm3.service
/usr/share/gdm/dconf/90-debian-settings
/usr/share/doc/gdm3/NEWS.gz
/usr/share/doc/gdm3/changelog.Debian.gz
/usr/share/doc/gdm3/README
/usr/share/doc/gdm3/AUTHORS
