# lib/dpkg/info/geoclue-2.0.list
/.
/usr
/usr/share
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.GeoClue2.service
/usr/share/dbus-1/interfaces
/usr/share/dbus-1/interfaces/org.freedesktop.GeoClue2.Manager.xml
/usr/share/dbus-1/interfaces/org.freedesktop.GeoClue2.Location.xml
/usr/share/dbus-1/interfaces/org.freedesktop.GeoClue2.xml
/usr/share/dbus-1/interfaces/org.freedesktop.GeoClue2.Agent.xml
/usr/share/dbus-1/interfaces/org.freedesktop.GeoClue2.Client.xml
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/geoclue-2.0
/usr/share/polkit-1
/usr/share/polkit-1/rules.d
/usr/share/polkit-1/rules.d/geoclue-2.0.rules
/usr/share/doc
/usr/share/doc/geoclue-2.0
/usr/share/doc/geoclue-2.0/copyright
/usr/share/applications
/usr/share/applications/geoclue-where-am-i.desktop
/usr/lib
/usr/lib/x86_64-linux-gnu
/usr/lib/x86_64-linux-gnu/pkgconfig
/usr/lib/x86_64-linux-gnu/pkgconfig/geoclue-2.0.pc
/usr/lib/geoclue-2.0
/usr/lib/geoclue-2.0/demos
/usr/lib/geoclue-2.0/demos/where-am-i
/usr/lib/geoclue-2.0/geoclue
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/geoclue.service
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/org.freedesktop.GeoClue2.Agent.conf
/etc/dbus-1/system.d/org.freedesktop.GeoClue2.conf
/etc/geoclue
/etc/geoclue/geoclue.conf
/var
/var/lib
/var/lib/polkit-1
/var/lib/polkit-1/localauthority
/var/lib/polkit-1/localauthority/10-vendor.d
/var/lib/polkit-1/localauthority/10-vendor.d/geoclue-2.0.pkla
/var/lib/geoclue
/usr/share/doc/geoclue-2.0/changelog.Debian.gz
