# lib/dpkg/info/iio-sensor-proxy.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/iio-sensor-proxy.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/40-iio-sensor-proxy.rules
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/net.hadess.SensorProxy.conf
/usr
/usr/share
/usr/share/doc
/usr/share/doc/iio-sensor-proxy
/usr/share/doc/iio-sensor-proxy/copyright
/usr/share/doc/iio-sensor-proxy/changelog.Debian.gz
/usr/sbin
/usr/sbin/iio-sensor-proxy
/usr/bin
/usr/bin/monitor-sensor
