# lib/dpkg/info/irqbalance.conffiles
/etc/init.d/irqbalance
/etc/init/irqbalance.conf
