# lib/dpkg/info/irqbalance.list
/.
/usr
/usr/share
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/irqbalance.1.gz
/usr/share/doc
/usr/share/doc/irqbalance
/usr/share/doc/irqbalance/copyright
/usr/share/doc/irqbalance/changelog.Debian.gz
/usr/sbin
/usr/sbin/irqbalance
/etc
/etc/init
/etc/init/irqbalance.conf
/etc/init.d
/etc/init.d/irqbalance
/etc/default
