# lib/dpkg/info/realmd.list
/.
/usr
/usr/share
/usr/share/locale
/usr/share/locale/ro/LC_MESSAGES/realmd.mo
/usr/share/locale/as/LC_MESSAGES/realmd.mo
/usr/share/locale/th/LC_MESSAGES/realmd.mo
/usr/share/locale/ms/LC_MESSAGES/realmd.mo
/usr/share/locale/sl/LC_MESSAGES/realmd.mo
/usr/share/locale/it_IT
/usr/share/locale/it_IT/LC_MESSAGES
/usr/share/locale/it_IT/LC_MESSAGES/realmd.mo
/usr/share/locale/vi/LC_MESSAGES/realmd.mo
/usr/share/locale/pl/LC_MESSAGES/realmd.mo
/usr/share/locale/ga/LC_MESSAGES/realmd.mo
/usr/share/locale/ta/LC_MESSAGES/realmd.mo
/usr/share/locale/te/LC_MESSAGES/realmd.mo
/usr/share/locale/id/LC_MESSAGES/realmd.mo
/usr/share/locale/ka/LC_MESSAGES/realmd.mo
/usr/share/locale/nl/LC_MESSAGES/realmd.mo
/usr/share/locale/sv/LC_MESSAGES/realmd.mo
/usr/share/locale/he/LC_MESSAGES/realmd.mo
/usr/share/locale/nb/LC_MESSAGES/realmd.mo
/usr/share/locale/fi/LC_MESSAGES/realmd.mo
/usr/share/locale/mr/LC_MESSAGES/realmd.mo
/usr/share/locale/az/LC_MESSAGES/realmd.mo
/usr/share/locale/zh_TW/LC_MESSAGES/realmd.mo
/usr/share/locale/ko/LC_MESSAGES/realmd.mo
/usr/share/locale/ca/LC_MESSAGES/realmd.mo
/usr/share/locale/kn/LC_MESSAGES/realmd.mo
/usr/share/locale/eu/LC_MESSAGES/realmd.mo
/usr/share/locale/ar/LC_MESSAGES/realmd.mo
/usr/share/locale/fo
/usr/share/locale/fo/LC_MESSAGES
/usr/share/locale/fo/LC_MESSAGES/realmd.mo
/usr/share/locale/ml/LC_MESSAGES/realmd.mo
/usr/share/locale/ca@valencia/LC_MESSAGES/realmd.mo
/usr/share/locale/pa/LC_MESSAGES/realmd.mo
/usr/share/locale/da/LC_MESSAGES/realmd.mo
/usr/share/locale/sk/LC_MESSAGES/realmd.mo
/usr/share/locale/lt/LC_MESSAGES/realmd.mo
/usr/share/locale/ja/LC_MESSAGES/realmd.mo
/usr/share/locale/el/LC_MESSAGES/realmd.mo
/usr/share/locale/de/LC_MESSAGES/realmd.mo
/usr/share/locale/gl/LC_MESSAGES/realmd.mo
/usr/share/locale/bn_IN/LC_MESSAGES/realmd.mo
/usr/share/locale/uk/LC_MESSAGES/realmd.mo
/usr/share/locale/sq/LC_MESSAGES/realmd.mo
/usr/share/locale/eo/LC_MESSAGES/realmd.mo
/usr/share/locale/kk/LC_MESSAGES/realmd.mo
/usr/share/locale/pt/LC_MESSAGES/realmd.mo
/usr/share/locale/nn/LC_MESSAGES/realmd.mo
/usr/share/locale/bg/LC_MESSAGES/realmd.mo
/usr/share/locale/fr/LC_MESSAGES/realmd.mo
/usr/share/locale/it/LC_MESSAGES/realmd.mo
/usr/share/locale/lv/LC_MESSAGES/realmd.mo
/usr/share/locale/hr/LC_MESSAGES/realmd.mo
/usr/share/locale/hu/LC_MESSAGES/realmd.mo
/usr/share/locale/zh_HK/LC_MESSAGES/realmd.mo
/usr/share/locale/wa/LC_MESSAGES/realmd.mo
/usr/share/locale/en_GB
/usr/share/locale/en_GB/LC_MESSAGES
/usr/share/locale/en_GB/LC_MESSAGES/realmd.mo
/usr/share/locale/ia/LC_MESSAGES/realmd.mo
/usr/share/locale/es/LC_MESSAGES/realmd.mo
/usr/share/locale/ru/LC_MESSAGES/realmd.mo
/usr/share/locale/cs/LC_MESSAGES/realmd.mo
/usr/share/locale/gu/LC_MESSAGES/realmd.mo
/usr/share/locale/or/LC_MESSAGES/realmd.mo
/usr/share/locale/pt_BR/LC_MESSAGES/realmd.mo
/usr/share/locale/zh_CN/LC_MESSAGES/realmd.mo
/usr/share/locale/cy/LC_MESSAGES/realmd.mo
/usr/share/locale/es_CL
/usr/share/locale/es_CL/LC_MESSAGES
/usr/share/locale/es_CL/LC_MESSAGES/realmd.mo
/usr/share/locale/hi/LC_MESSAGES/realmd.mo
/usr/share/locale/sr@latin/LC_MESSAGES/realmd.mo
/usr/share/locale/fa/LC_MESSAGES/realmd.mo
/usr/share/locale/tr/LC_MESSAGES/realmd.mo
/usr/share/doc
/usr/share/doc/realmd
/usr/share/doc/realmd/html
/usr/share/doc/realmd/html/guide-active-directory-join.html
/usr/share/doc/realmd/html/left.png
/usr/share/doc/realmd/html/gdbus-org.freedesktop.realmd.Kerberos.html
/usr/share/doc/realmd/html/dbus-interface-reference.html
/usr/share/doc/realmd/html/guide-ipa.html
/usr/share/doc/realmd/html/dbus-interface-raw.html
/usr/share/doc/realmd/html/guide-kerberos.html
/usr/share/doc/realmd/html/guide-active-directory-client.html
/usr/share/doc/realmd/html/gdbus-org.freedesktop.realmd.Provider.html
/usr/share/doc/realmd/html/guide-ipa-join.html
/usr/share/doc/realmd/html/guide-integration.html
/usr/share/doc/realmd/html/gdbus-org.freedesktop.realmd.Service.html
/usr/share/doc/realmd/html/style.css
/usr/share/doc/realmd/html/right.png
/usr/share/doc/realmd/html/realmd-conf.html
/usr/share/doc/realmd/html/development.html
/usr/share/doc/realmd/html/up.png
/usr/share/doc/realmd/html/guide-ipa-client.html
/usr/share/doc/realmd/html/gdbus-org.freedesktop.realmd.Realm.html
/usr/share/doc/realmd/html/index.html
/usr/share/doc/realmd/html/guide.html
/usr/share/doc/realmd/html/guide-active-directory.html
/usr/share/doc/realmd/html/gdbus-org.freedesktop.realmd.KerberosMembership.html
/usr/share/doc/realmd/html/realm.html
/usr/share/doc/realmd/html/gtk-doc.css
/usr/share/doc/realmd/html/guide-ipa-permit.html
/usr/share/doc/realmd/html/realm-manual.html
/usr/share/doc/realmd/html/index.sgml
/usr/share/doc/realmd/html/home.png
/usr/share/doc/realmd/html/guide-active-directory-permit.html
/usr/share/doc/realmd/html/realmd-docs.proc
/usr/share/doc/realmd/README
/usr/share/doc/realmd/copyright
/usr/share/doc/realmd/changelog.Debian.gz
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.freedesktop.realmd.policy
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/realmd.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/realm.8.gz
/usr/share/doc-base
/usr/share/doc-base/realmd
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.realmd.service
/usr/sbin
/usr/sbin/realm
/usr/lib
/usr/lib/realmd
/usr/lib/realmd/realmd-defaults.conf
/usr/lib/realmd/realmd
/usr/lib/realmd/realmd-distro.conf
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/org.freedesktop.realmd.conf
/var
/var/cache
/var/cache/realmd
/var/lib
/var/lib/realmd
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/realmd.service
