# lib/dpkg/info/resolvconf.list
/.
/lib
/lib/resolvconf
/lib/resolvconf/list-records
/lib/systemd
/lib/systemd/system
/lib/systemd/system/resolvconf.service
/etc
/etc/network
/etc/network/if-down.d
/etc/network/if-down.d/resolvconf
/etc/network/if-up.d
/etc/network/if-up.d/000resolvconf
/etc/init
/etc/init/resolvconf.conf
/etc/dhcp
/etc/dhcp/dhclient-enter-hooks.d
/etc/dhcp/dhclient-enter-hooks.d/resolvconf
/etc/resolvconf
/etc/resolvconf/interface-order
/etc/resolvconf/resolv.conf.d
/etc/resolvconf/resolv.conf.d/base
/etc/resolvconf/resolv.conf.d/head
/etc/resolvconf/update.d
/etc/resolvconf/update.d/libc
/etc/ppp
/etc/ppp/ip-down.d
/etc/ppp/ip-down.d/000resolvconf
/etc/ppp/ip-up.d
/etc/ppp/ip-up.d/000resolvconf
/etc/init.d
/etc/init.d/resolvconf
/var
/var/lib
/var/lib/resolvconf
/usr
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/resolvconf.8.gz
/usr/share/man/man5
/usr/share/man/man5/interface-order.5.gz
/usr/share/resolvconf
/usr/share/resolvconf/dump-debug-info
/usr/share/doc
/usr/share/doc/resolvconf
/usr/share/doc/resolvconf/copyright
/usr/share/doc/resolvconf/README.gz
/usr/share/doc/resolvconf/changelog.gz
/usr/share/doc/resolvconf/NEWS.Debian.gz
/usr/share/doc/resolvconf/resolvconf-update-bind
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/resolvconf
/sbin
/sbin/resolvconf
