# lib/dpkg/info/snapd.list
/.
/var
/var/snap
/var/lib
/var/lib/snapd
/var/lib/snapd/firstboot
/var/lib/snapd/snaps
/var/lib/snapd/snaps/partial
/var/lib/snapd/auto-import
/var/lib/snapd/void
/var/lib/snapd/environment
/var/lib/snapd/lib
/var/lib/snapd/lib/gl
/var/lib/snapd/desktop
/snap
/usr
/usr/bin
/usr/bin/snap
/usr/bin/snapctl
/usr/lib
/usr/lib/snapd
/usr/lib/snapd/snapd
/usr/lib/snapd/snap-exec
/usr/lib/snapd/snapd.core-fixup.sh
/usr/lib/snapd/system-shutdown
/usr/lib/snapd/snap-confine
/usr/lib/snapd/snap-discard-ns
/usr/lib/snapd/etelpmoc.sh
/usr/lib/snapd/snap-repair
/usr/lib/snapd/complete.sh
/usr/lib/snapd/snap-update-ns
/usr/lib/snapd/info
/usr/lib/snapd/snap-seccomp
/usr/share
/usr/share/doc
/usr/share/doc/snapd
/usr/share/doc/snapd/changelog.gz
/usr/share/doc/snapd/copyright
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/snap-confine.5.gz
/usr/share/man/man5/snap-discard-ns.5.gz
/usr/share/man/man1
/usr/share/man/man1/snap.1.gz
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/snap
/lib
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/66-snapd-autoimport.rules
/lib/udev/rules.d/80-snappy-assign.rules
/lib/udev/snappy-app-dev
/lib/systemd
/lib/systemd/system
/lib/systemd/system/snapd.autoimport.service
/lib/systemd/system/snapd.refresh.service
/lib/systemd/system/snapd.refresh.timer
/lib/systemd/system/snap-repair.service
/lib/systemd/system/snap-repair.timer
/lib/systemd/system/snapd.core-fixup.service
/lib/systemd/system/snapd.socket
/lib/systemd/system/snapd.system-shutdown.service
/lib/systemd/system/snapd.service
/etc
/etc/X11
/etc/X11/Xsession.d
/etc/X11/Xsession.d/65snappy
/etc/apparmor.d
/etc/apparmor.d/usr.lib.snapd.snap-confine.real
/etc/profile.d
/etc/profile.d/apps-bin-path.sh
/usr/bin/ubuntu-core-launcher
