# lib/dpkg/info/system-config-printer-udev.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/udev-configure-printer@.service
/lib/udev
/lib/udev/udev-configure-printer
/lib/udev/rules.d
/lib/udev/rules.d/70-printers.rules
/lib/udev/udev-add-printer
/usr
/usr/share
/usr/share/doc
/usr/share/doc/system-config-printer-udev
/usr/share/doc/system-config-printer-udev/AUTHORS
/usr/share/doc/system-config-printer-udev/copyright
/usr/share/doc/system-config-printer-udev/changelog.Debian.gz
