# lib/dpkg/info/tor.conffiles
/etc/apparmor.d/abstractions/tor
/etc/apparmor.d/system_tor
/etc/cron.weekly/tor
/etc/default/tor
/etc/init.d/tor
/etc/logrotate.d/tor
/etc/tor/torrc
