# lib/dpkg/info/tor.list
/.
/usr
/usr/sbin
/usr/sbin/tor-instance-create
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/tor
/usr/share/tor
/usr/share/tor/tor-service-defaults-torrc-instances
/usr/share/tor/tor-service-defaults-torrc
/usr/share/man
/usr/share/man/man5
/usr/share/man/man8
/usr/share/man/man8/tor-instance-create.8.gz
/usr/share/man/man1
/usr/share/man/man1/torify.1.gz
/usr/share/man/man1/tor.1.gz
/usr/share/man/man1/tor-resolve.1.gz
/usr/share/man/man1/tor-gencert.1.gz
/usr/share/doc
/usr/share/doc/tor
/usr/share/doc/tor/changelog.Debian.gz
/usr/share/doc/tor/NEWS.Debian.gz
/usr/share/doc/tor/copyright
/usr/share/doc/tor/README.Debian
/usr/share/doc/tor/tor-exit-notice.html
/usr/share/doc/tor/torify.html
/usr/share/doc/tor/tor-resolve.html
/usr/share/doc/tor/tor-gencert.html
/usr/share/doc/tor/tor.html
/usr/bin
/usr/bin/tor
/usr/bin/tor-resolve
/usr/bin/tor-gencert
/usr/bin/torify
/lib
/lib/systemd
/lib/systemd/system-generators
/lib/systemd/system-generators/tor-generator
/lib/systemd/system
/lib/systemd/system/tor@default.service
/lib/systemd/system/tor@.service
/lib/systemd/system/tor.service
/etc
/etc/logrotate.d
/etc/logrotate.d/tor
/etc/init.d
/etc/init.d/tor
/etc/default
/etc/default/tor
/etc/cron.weekly
/etc/cron.weekly/tor
/etc/tor
/etc/tor/torrc
/etc/apparmor.d
/etc/apparmor.d/system_tor
/etc/apparmor.d/abstractions
/etc/apparmor.d/abstractions/tor
/usr/sbin/tor
/usr/share/man/man5/torrc.5.gz
