# lib/dpkg/info/ureadahead.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/ureadahead-stop.service
/lib/systemd/system/ureadahead-stop.timer
/lib/systemd/system/ureadahead.service
/etc
/etc/init
/etc/init/ureadahead-other.conf
/etc/init/ureadahead.conf
/sbin
/sbin/ureadahead
/var
/var/lib
/var/lib/ureadahead
/var/lib/ureadahead/debugfs
/usr
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/ureadahead.8.gz
/usr/share/doc
/usr/share/doc/ureadahead
/usr/share/doc/ureadahead/changelog.Debian.gz
/usr/share/doc/ureadahead/copyright
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/ureadahead.py
