# lib/dpkg/info/uuid-runtime.list
/.
/usr
/usr/bin
/usr/bin/uuidgen
/usr/share
/usr/share/doc
/usr/share/doc/uuid-runtime
/usr/share/doc/uuid-runtime/changelog.Debian.gz
/usr/share/doc/uuid-runtime/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/uuidgen.1.gz
/usr/share/man/man8
/usr/share/man/man8/uuidd.8.gz
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/uuidd
/usr/share/bash-completion/completions/uuidgen
/usr/sbin
/usr/sbin/uuidd
/etc
/etc/init.d
/etc/init.d/uuidd
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/uuidd.socket
/lib/systemd/system/uuidd.service
