# lib/dpkg/info/virtualbox.list
/.
/usr
/usr/share
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_virtualbox.py
/usr/share/virtualbox
/usr/share/virtualbox/VBox.sh
/usr/share/virtualbox/VBoxSysInfo.sh
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/virtualbox
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/vboxwebsrv.1.gz
/usr/share/man/man1/VBoxSDL.1.gz
/usr/share/man/man1/VBoxHeadless.1.gz
/usr/share/man/man1/VBoxManage.1.gz
/usr/share/python
/usr/share/python/runtime.d
/usr/share/python/runtime.d/virtualbox.rtupdate
/usr/share/doc
/usr/share/doc/virtualbox
/usr/share/doc/virtualbox/README.Debian.gz
/usr/share/doc/virtualbox/changelog.Debian.gz
/usr/share/doc/virtualbox/UserManual.pdf
/usr/share/doc/virtualbox/NEWS.Debian.gz
/usr/share/doc/virtualbox/copyright
/usr/share/doc/virtualbox/README.Debian.html
/usr/share/doc/virtualbox/wiki.html
/usr/bin
/usr/lib
/usr/lib/virtualbox
/usr/lib/virtualbox/libvboxjxpcom.so
/usr/lib/virtualbox/VBoxExtPackHelperApp
/usr/lib/virtualbox/sdk
/usr/lib/virtualbox/sdk/bindings
/usr/lib/virtualbox/sdk/bindings/xpcom
/usr/lib/virtualbox/sdk/bindings/xpcom/python
/usr/lib/virtualbox/sdk/bindings/xpcom/python/xpcom
/usr/lib/virtualbox/sdk/bindings/xpcom/python/xpcom/vboxxpcom.py
/usr/lib/virtualbox/sdk/bindings/xpcom/python/xpcom/server
/usr/lib/virtualbox/sdk/bindings/xpcom/python/xpcom/server/factory.py
/usr/lib/virtualbox/sdk/bindings/xpcom/python/xpcom/server/module.py
/usr/lib/virtualbox/sdk/bindings/xpcom/python/xpcom/server/enumerator.py
/usr/lib/virtualbox/sdk/bindings/xpcom/python/xpcom/server/loader.py
/usr/lib/virtualbox/sdk/bindings/xpcom/python/xpcom/server/policy.py
/usr/lib/virtualbox/sdk/bindings/xpcom/python/xpcom/server/__init__.py
/usr/lib/virtualbox/sdk/bindings/xpcom/python/xpcom/components.py
/usr/lib/virtualbox/sdk/bindings/xpcom/python/xpcom/primitives.py
/usr/lib/virtualbox/sdk/bindings/xpcom/python/xpcom/xpt.py
/usr/lib/virtualbox/sdk/bindings/xpcom/python/xpcom/client
/usr/lib/virtualbox/sdk/bindings/xpcom/python/xpcom/client/__init__.py
/usr/lib/virtualbox/sdk/bindings/xpcom/python/xpcom/xpcom_consts.py
/usr/lib/virtualbox/sdk/bindings/xpcom/python/xpcom/nsError.py
/usr/lib/virtualbox/sdk/bindings/xpcom/python/xpcom/file.py
/usr/lib/virtualbox/sdk/bindings/xpcom/python/xpcom/__init__.py
/usr/lib/virtualbox/sdk/bindings/xpcom/java
/usr/lib/virtualbox/sdk/bindings/xpcom/java/vboxjxpcom.jar
/usr/lib/virtualbox/VBoxXPCOM.so
/usr/lib/virtualbox/VBoxKeyboard.so
/usr/lib/virtualbox/VMMRC.rc
/usr/lib/virtualbox/VBoxDD2RC.rc
/usr/lib/virtualbox/VBoxHeadless
/usr/lib/virtualbox/VBoxGuestControlSvc.so
/usr/lib/virtualbox/VMMR0.r0
/usr/lib/virtualbox/VBoxDDRC.rc
/usr/lib/virtualbox/VBoxNetDHCP
/usr/lib/virtualbox/VBoxPython2_7.so
/usr/lib/virtualbox/VBoxPython.so
/usr/lib/virtualbox/VBoxEFI64.fd
/usr/lib/virtualbox/VBoxSDL.so
/usr/lib/virtualbox/VBoxGuestPropSvc.so
/usr/lib/virtualbox/VBoxOGLrenderspu.so
/usr/lib/virtualbox/VBoxOGLhostcrutil.so
/usr/lib/virtualbox/VBoxNetNAT.so
/usr/lib/virtualbox/VBoxBugReport
/usr/lib/virtualbox/VBoxNetAdpCtl
/usr/lib/virtualbox/VBoxVMMPreload.so
/usr/lib/virtualbox/VBoxSharedClipboard.so
/usr/lib/virtualbox/VBoxXPCOMC.so
/usr/lib/virtualbox/webtest
/usr/lib/virtualbox/VBoxSharedFolders.so
/usr/lib/virtualbox/VBoxBalloonCtrl
/usr/lib/virtualbox/VBoxHostChannel.so
/usr/lib/virtualbox/VBoxDragAndDropSvc.so
/usr/lib/virtualbox/VBoxSVC
/usr/lib/virtualbox/components
/usr/lib/virtualbox/components/VirtualBox_XPCOM.xpt
/usr/lib/virtualbox/components/VBoxXPCOMBase.xpt
/usr/lib/virtualbox/components/VBoxSVCM.so
/usr/lib/virtualbox/components/VBoxC.so
/usr/lib/virtualbox/components/VBoxXPCOMIPCC.so
/usr/lib/virtualbox/VBoxOGLhosterrorspu.so
/usr/lib/virtualbox/VBoxAuthSimple.so
/usr/lib/virtualbox/VBoxDDU.so
/usr/lib/virtualbox/VBoxAuth.so
/usr/lib/virtualbox/VBoxTestOGL
/usr/lib/virtualbox/vboxshell.py
/usr/lib/virtualbox/VBoxREM.so
/usr/lib/virtualbox/VBoxDD2.so
/usr/lib/virtualbox/VBoxAutostart
/usr/lib/virtualbox/VBoxHeadless.so
/usr/lib/virtualbox/VBoxVMM.so
/usr/lib/virtualbox/VBoxDD.so
/usr/lib/virtualbox/VBoxSharedCrOpenGL.so
/usr/lib/virtualbox/ExtensionPacks
/usr/lib/virtualbox/ExtensionPacks/VNC
/usr/lib/virtualbox/ExtensionPacks/VNC/ExtPack-license.rtf
/usr/lib/virtualbox/ExtensionPacks/VNC/ExtPack.xml
/usr/lib/virtualbox/ExtensionPacks/VNC/ExtPack-license.txt
/usr/lib/virtualbox/ExtensionPacks/VNC/ExtPack-license.html
/usr/lib/virtualbox/ExtensionPacks/VNC/linux.amd64
/usr/lib/virtualbox/ExtensionPacks/VNC/linux.amd64/VBoxVNCMain.so
/usr/lib/virtualbox/ExtensionPacks/VNC/linux.amd64/VBoxVNC.so
/usr/lib/virtualbox/VBoxXPCOMIPCD
/usr/lib/virtualbox/VBoxDDR0.r0
/usr/lib/virtualbox/vboxwebsrv
/usr/lib/virtualbox/VBoxNetDHCP.so
/usr/lib/virtualbox/VBoxManage
/usr/lib/virtualbox/VBoxRT.so
/usr/lib/virtualbox/VBoxDD2R0.r0
/usr/lib/virtualbox/VBoxEFI32.fd
/usr/lib/virtualbox/VBoxSDL
/usr/lib/virtualbox/vbox-img
/usr/lib/python2.7
/usr/lib/python2.7/dist-packages
/usr/lib/python2.7/dist-packages/vboxapi-1.0.egg-info
/usr/lib/python2.7/dist-packages/vboxapi
/usr/lib/python2.7/dist-packages/vboxapi/VirtualBox_constants.py
/usr/lib/python2.7/dist-packages/vboxapi/__init__.py
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/vboxweb.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/60-virtualbox.rules
/lib/udev/VBoxCreateUSBNode.sh
/etc
/etc/default
/etc/default/virtualbox
/etc/init.d
/etc/init.d/virtualbox
/usr/share/man/man1/vboxsdl.1.gz
/usr/share/man/man1/vboxmanage.1.gz
/usr/share/man/man1/vboxheadless.1.gz
/usr/bin/VBoxHeadless
/usr/bin/vboxballoonctrl
/usr/bin/vboxsdl
/usr/bin/vboxheadless
/usr/bin/VBoxBalloonCtrl
/usr/bin/vboxmanage
/usr/bin/vboxwebsrv
/usr/bin/VBoxManage
/usr/bin/VBoxSDL
/usr/bin/vbox-img
/usr/lib/virtualbox/components/VBoxXPCOM.so
/usr/lib/virtualbox/components/VBoxDDU.so
/usr/lib/virtualbox/components/VBoxREM.so
/usr/lib/virtualbox/components/VBoxVMM.so
/usr/lib/virtualbox/components/VBoxRT.so
