# lib/dpkg/info/virtualbox.conffiles
/etc/default/virtualbox
/etc/init.d/virtualbox
