# etc/init/network-manager.conf
# network-manager - network connection manager
#
# The Network Manager daemon manages the system's network connections,
# automatically switching between the best available.

description	"network connection manager"

start on (local-filesystems
	  and started dbus
	  and static-network-up)
stop on stopping dbus

expect fork
respawn

script
	# set $LANG so that messages appearing on the GUI will be translated. See LP: 875017
	if [ -r /etc/default/locale ]; then
		. /etc/default/locale
		export LANG LANGUAGE LC_MESSAGES LC_ALL
	fi

	exec NetworkManager
end script
