# etc/init/usb-modeswitch-upstart.conf
start on usb-modeswitch-upstart
task
env TMPDIR=/run
script
	exec /usr/sbin/usb_modeswitch_dispatcher --switch-upstart $UMS_PARAM
end script
